# cran-notifier.timer
[Unit]
Description=Run CRAN incoming checker every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
